# supervisord.conf
[supervisord]
nodaemon=true
user=root

# 1. Jalankan Nginx Unit (Web Server & PHP)
[program:unitd]
command=/usr/sbin/unitd --no-daemon --control unix:/var/run/control.unit.sock
user=root
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=1

# 2. Startup Tasks Runner & Config Loader
[program:unit-config-loader]
command=bash -c "echo '🚀 Starting Deployment Tasks...'; chown -R unit:unit /var/www/html; chmod -R 775 /var/www/html/storage /var/www/html/bootstrap/cache; echo '🔗 Linking Storage...'; php artisan storage:link; echo '📦 Migrating Database...'; php artisan migrate --force; echo '⚡ Caching Config...'; php artisan optimize; echo '⏳ Waiting for Unit...'; sleep 5; curl -X PUT --data-binary @/docker-entrypoint.d/unit.json --unix-socket /var/run/control.unit.sock http://localhost/config; echo '✅ System Ready!'"
autostart=true
autorestart=false
startretries=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

# 3. Laravel Queue Worker (Production)
[program:queue-worker]
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3 --max-time=3600
user=unit
directory=/var/www/html
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

# 4. MQTT → InfluxDB Bridge (replaces EMQX Rule Engine)
[program:mqtt-listener]
command=php /var/www/html/artisan mqtt:listen
user=unit
directory=/var/www/html
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30
startsecs=5
startretries=50
